# deploy/systemd/stocklean-weekly-rebalance-prepare.timer
[Unit]
Description=Run StockLean Weekly Rebalance Prepare

[Timer]
OnCalendar=Mon *-*-* 08:00:00
Persistent=true
Unit=stocklean-weekly-rebalance-prepare.service

[Install]
WantedBy=timers.target
