# deploy/systemd/stocklean-weekly-rebalance-execute.service
[Unit]
Description=StockLean Weekly Rebalance Execute
After=stocklean-backend.service

[Service]
Type=oneshot
WorkingDirectory=/app/stocklean
Environment=STOCKLEAN_ROOT=/app/stocklean
Environment=STOCKLEAN_BACKEND_URL=http://127.0.0.1:8021
Environment=STOCKLEAN_PROJECT_ID=18
ExecStart=/app/stocklean/scripts/weekly_rebalance.sh execute
